# .config/autostart/xinit.desktop
[Desktop Entry]
Comment=xinit
Exec=/bin/bash /home/steamedfish/.xinitrc
Path=/home/steamedfish/
StartupNotify=false
Terminal=true
Type=Application
